# chirp-impl/src/main/docker/root/etc/services.d/auth/run
#!/usr/bin/with-contenv sh

cd /opt/docker

bin/chirp-impl -Dplay.crypto.secret="asdfasdfasdfasdf"
