# front-end/docker/root/etc/services.d/auth/run
#!/usr/bin/with-contenv sh

cd /opt/docker

bin/front-end -Dplay.crypto.secret="asdfasdfasdfasdf"
